# ops/systemd/hs-source-scan.service
[Unit]
Description=HS Tracker Sponsor Source Scan (PR + Pipeline Pages)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=<CODEX_REPO>
Environment=HS_TRACKER_DB_PATH=<CODEX_REPO>/data/hs_tracker.db
ExecStart=<PYTHON_BIN> -m hs_tracker.jobs.run_source_scan --config data/source_configs/sponsor_sources.json
User=<SERVICE_USER>
Group=<SERVICE_GROUP>
